# stow/apps/.crontab
# composer auto-update every month
@monthly /opt/homebrew/bin/composer self-update >> /dev/null 2>&1

# npm update version every month
@monthly /opt/homebrew/bin/npm install --global npm@latest >> /dev/null 2>&1
